# openmptcprouter/files/etc/uci-defaults/1950-omr-diag
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	set luci.diag=internal
	set luci.diag.ping="zjzwyl.com"
	set luci.diag.dns="zjzwyl.com"
	set luci.diag.route="zjzwyl.com"
	set luci.diag.iperf3="ping-ams1.online.net"
	set luci.diag.getip="ip.openmptcprouter.com"
	commit luci
EOF
rm -f /tmp/luci-indexcache

exit 0
